# s6-overlay/s6-rc.d/04-soapysdr/run
#!/command/with-contenv sh
echo "Starting soapysdr"
exec /usr/local/bin/start-soapysdr.sh
